# conversations/scripts/com.humboldtscabinet.claude-conversation-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.humboldtscabinet.claude-conversation-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <!-- -l makes sure Homebrew's git/python are on PATH when launchd runs us -->
        <string>"$HOME/Developer/organizing-life-services-ai/conversations/scripts/backup_conversations.sh"</string>
    </array>

    <!-- Run daily at 2:15am local time. Change Hour/Minute if you prefer. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <!-- Manual trigger: touching this file fires the backup. Works from
         any shell, and from Claude (the file is inside the workspace folder
         that the sandbox can write to). -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/robertporter/Developer/organizing-life-services-ai/conversations/.trigger</string>
    </array>

    <!-- If the Mac was asleep at 2:15, run as soon as it wakes -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Log output to the same folder as the script -->
    <key>StandardOutPath</key>
    <string>/tmp/claude-conversation-backup.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-conversation-backup.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
